# workspace/ct-dmz-proxy/gate-monitoring/ctdmz-gate-metrics.timer
[Unit]
Description=Run ctdmz-gate-metrics every 30s (refresh gate-health textfile)

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=ctdmz-gate-metrics.service

[Install]
WantedBy=timers.target
